# etc/supervisor/conf.d/httpd.conf
[program:httpd]
directory=/data/modules/httpd
command=/data/modules/httpd/bin/httpd --smp 4
autostart=false
autorestart=true
